# logCPUTemp.service
[Unit]
Description=Log CPU temperature

[Service]
ExecStart=[pythonのパス] [main.pyのパス]
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target